# server/deploy/weedit-signaling.service
[Unit]
Description=WeEdit WebRTC signaling server
Documentation=https://github.com/minecartchris/WeEdit
After=network.target

[Service]
Type=simple
# Runs as an unprivileged user (create with: sudo useradd --system --no-create-home --shell /usr/sbin/nologin weedit)
User=weedit
Group=weedit
WorkingDirectory=/opt/weedit-signaling
ExecStart=/usr/bin/node /opt/weedit-signaling/signaling.mjs
Environment=PORT=4444
# 0.0.0.0 = reachable on the LAN / via a tunnel. Set to 127.0.0.1 if you want
# only a local nginx reverse proxy to be able to reach it.
Environment=HOST=0.0.0.0
Restart=on-failure
RestartSec=2

# Hardening — the server needs nothing but a loopback socket.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=
ProtectKernelTunables=true
ProtectControlGroups=true
# AF_UNIX/AF_NETLINK included because Node/libuv can need them at startup
# (interface enumeration, internal pipes); without them the process can fail
# to start under the sandbox.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK

[Install]
WantedBy=multi-user.target
